# backend/scripts/ksc-mcube-controller.service
[Unit]
Description=KylinSec security magic Cube controller service

[Service]
ExecStart=bash -c "/usr/bin/ksc-mcube-server $(cat /etc/ksc-mcube/server.flags) controller"

Restart=on-failure

RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target